# installers/macos/root/Library/LaunchAgents/com.anava.local_connector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.anava.local_connector</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/Applications/AnavaLocalConnector/local-connector</string>
        <string>--proxy-service</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/anava-proxy-service.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/anava-proxy-service-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
